# crontab/cronbuild
#!/bin/bash

pkill -f 'cronbuilder'
jps | cut -d' ' -f 1 xargs kill
$(dirname $0)/cronbuilder $@