# crontab/crondownload
#!/bin/bash

pkill -f 'crondownloader'
$(dirname $0)/crondownloader