# crontab/cronbuild
#!/bin/bash

pkill -f 'cronbuilder'
./cronbuilder