# crontab/cronbuild
#!/bin/bash

pkill -f 'cronbuilder'
$(dirname $0)/cronbuilder $@